# osx_launchctl/runit/org.smarden.runit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.smarden.runit</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>MachServices</key>
  <dict>
    <key>org.smarden.runit</key>
    <true/>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
	<array>
    <string>/usr/local/bin/runsvdir</string>
    <string>/service</string>
	</array>
</dict>
</plist>
